# workshops/dash-observability-pipelines/systemd/storedog-gen.service
[Unit]
Description=Generates synthetic logs from an HTTP service.

[Service]
ExecStart=/root/lab/bin/storedog-gen
Restart=always

[Install]
WantedBy=multi-user.target